# PolluxConfig/debian/postrm
#!/bin/sh
#systemctl disable pollux_config.service

[ -f "$D/etc/lighttpd.conf" ] && sed -i '/^\(.*pollux_config.*\)$/d' $D/etc/lighttpd.conf
