# PolluxContribIpks/lighttpd-systemd/debian/postinst
#!/bin/sh
OPTS=""

if [ -n "$D" ]; then
    OPTS="--root=$D"
fi

systemctl $OPTS enable lighttpd.service

if [ -z "$D" ]; then
    systemctl start lighttpd.service
fi

