# PolluxContribIpks/lighttpd-systemd/debian/postrm
#!/bin/sh
systemctl disable lighttpd.service
