# PolluxContribIpks/pollux-full/etc/cron.hourly/keep_time_in_sync
#!/bin/sh

/usr/bin/ntpdate -b -s -u pool.ntp.org
